# examples/tns.service
[Unit]
Description=Token Never Sleeps runner
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /absolute/path/to/token-never-sleeps/scripts/tns_runner.py run --config /absolute/path/to/tns_config.json
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
